# stage3/03-crankshaft-base/files/etc/systemd/system/disconnect.timer
[Unit]
 Description=Crankshaft Disconnect Timer

[Timer]
 OnActiveSec=120s
 AccuracySec=1s

[Install]
 WantedBy=timers.target
