# stage3/03-crankshaft-base/files/etc/systemd/system/openauto.service
[Unit]
 Description=Crankshaft Android Auto Headunit
 After=usbdetect.service usbrestore.service pulseaudio.service
 Before=timerstart.service

[Service]
 Type=simple
 User=pi
 WorkingDirectory=/home/pi
 ExecStart=/opt/crankshaft/service_openauto.sh
 Restart=on-failure
 StartLimitInterval=60
 RestartSec=2
 StartLimitBurst=3

[Install]
 WantedBy=multi-user.target
