# stage2/03-crankshaft-base/files/etc/systemd/system/alsastaterestore.service
[Unit]
 Description=Crankshaft alsa state restore
 After=sysinit.target

[Service]
 Type=oneshot
 ExecStart=/usr/local/bin/crankshaft audio alsastate restore
