# stage2/03-crankshaft-base/files/etc/systemd/system/custombrightness.service
[Unit]
 Description=Crankshaft Custom Brightness Exec Service

[Service]
 Type=simple
 ExecStart=/opt/crankshaft/service_custombrightness.sh
 RemainAfterExit=yes
