# stage2/03-crankshaft-base/files/etc/systemd/system/daymode.timer
[Unit]
Description=Crankshaft Day Mode Timer

[Timer]
OnCalendar=*-*-* 8..17:*:00

[Install]
WantedBy=timers.target
