# stage2/03-crankshaft-base/files/etc/systemd/system/daymode.timer
[Unit]
 Description=Crankshaft Day Mode Timer

[Timer]
 OnCalendar=*-*-* 8..19:*:00
 AccuracySec=5s

[Install]
 WantedBy=timers.target
