# stage2/03-crankshaft-base/files/etc/systemd/system/daymode.service
[Unit]
 Description=Crankshaft Day Mode Service

[Service]
 Type=simple
 ExecStart=/opt/crankshaft/service_daynight.sh day
