# stage2/03-crankshaft-base/files/etc/systemd/system/kodimonitor.service
[Unit]
 Description=Crankshaft Kodi Monitor Service

[Service]
 Type=simple
 ExecStart=/opt/crankshaft/service_kodimonitor.sh
 RemainAfterExit=yes
