# stage2/03-crankshaft-base/files/etc/systemd/system/nightmode.service
[Unit]
 Description=Crankshaft Night Mode Service

[Service]
 Type=simple
 ExecStart=/opt/crankshaft/service_daynight.sh night

