# stage2/03-crankshaft-base/files/etc/systemd/system/shutdown.timer
[Unit]
 Description=Crankshaft Shutdown Timer

[Timer]
 OnActiveSec=60m
 AccuracySec=5s

[Install]
 WantedBy=timers.target
