# stage2/04-crankshaft-bluetooth/files/etc/systemd/system/btautoconnect.service
[Unit]
 Description=Crankshaft Bluetooth Auto Connect Service

[Service]
 Type=oneshot
 User=pi
 WorkingDirectory=/home/pi
 ExecStart=/opt/crankshaft/service_btautoconnect.sh
 RemainAfterExit=yes

[Install]
 WantedBy=multi-user.target
