# stage2/04-crankshaft-bluetooth/files/etc/systemd/system/csng-bluetooth.service
[Unit]
 Description=Crankshaft Bluetooth Service
 After=bluetooth.target dbus-org.bluez.service sound.target rc.local.service

[Service]
 Type=oneshot
 User=pi
 WorkingDirectory=/home/pi
 ExecStart=/opt/crankshaft/service_bluetooth.sh
 RemainAfterExit=yes

[Install]
 WantedBy=multi-user.target
