# precompiled/systemd_services/autoapp.service
[Unit]
Description=Android Auto Headunit

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi
ExecStart=/usr/local/bin/autoapp
#ExecStopPost=/sbin/shutdown --halt now
Restart=on-abort

[Install]
WantedBy=multi-user.target

